# docker/config/php/supervisor/mediator-consume.conf
[program:mediator-consume]
command=php /var/www/artisan broker:consume:mediator 10 10
user=root
numprocs=1
startsecs=0
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d
